# templates/init/systemd.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          scaffold
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Runs Scaffold as a service.
# Description:       Runs Scaffold as a service.
### END INIT INFO

DAEMON="/opt/scaffold/bin/scaffold"
DAEMON_OPTS=""
DAEMON_USER="scaffold"
DAEMON_HOME="/opt/scaffold"
DAEMON_NAME="scaffold"
DAEMON_PID="/var/run/$DAEMON_NAME.pid"

. /lib/init/vars.sh
. /lib/lsb/init-functions

# Handle setup lifecycle events.
do_start() {
  log_daemon_msg "Starting $DAEMON_NAME"
  start-stop-daemon \
    --start \
    --background \
    --pidfile $DAEMON_PID \
    --make-pidfile \
    --user $DAEMON_USER \
    --chuid $DAEMON_USER \
    --startas $DAEMON \
    -- $DAEMON_OPTS
  log_end_msg $?
}

# Handle shutdown lifecycle events.
do_shutdown() {
  log_daemon_msg "Stopping $DAEMON_NAME"
  start-stop-daemon \
    --stop \
    --pidfile $DAEMON_PID \
    --retry 30
  log_end_msg $?
}

case "$1" in

  start)
    do_start
    ;;

  stop)
    do_shutdown
    ;;

  status)
    status_of_proc "$DAEMON_NAME" "$DAEMON" && exit 0 || exit $?
    ;;

  *)
    echo "Usage: $NAME {start|stop|status}" >&2
    exit 1
    ;;

esac
